# crontab
*/20 * * * * connectToVala
@reboot /home/dian/projects/dotfiles/scripts/batteryAlert/battery

